# deploy/systemd/mtx-dashboard.service
[Unit]
Description=MediaMTX Dashboard (Vite Preview)
After=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mtx/dashboard
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 5173
Restart=always
RestartSec=2
LimitNOFILE=262144

[Install]
WantedBy=multi-user.target
